# .config/systemd/user/graphical@.target
[Unit]
Description=Graphical interface using %i